# bin/jbdev.service
[Unit]
Description=Joe B.'s Home Automation Bus (dev)
After=network.target
Documentation=https://github.com/r1ptide64/jbe

[Service]
Restart=always
RestartSec=10
User=root
Group=root
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=jbdev
Environment=DEBUG=*
ExecStart=/usr/bin/node /home/joeb/dev/bin/www.js

[Install]
WantedBy=multi-user.target