# bin/jbprd.service
[Unit]
Description=Joe B.'s Home Automation Bus
After=network.target
Documentation=https://github.com/r1ptide64/jbe

[Service]
Restart=always
RestartSec=10
User=root
Group=root
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=jbprd
Environment=DEBUG=* NODE_ENV=production
ExecStart=/usr/bin/node /home/joeb/prd/bin/www.js

[Install]
WantedBy=multi-user.target